# crontab
*/5 * * * * /home/user/.local/bin/batterynotify.sh
